# src/assembly/bin/bard-crud-sample.service
[Unit]
Description=Bard CRUD Sample
Requires=network.target

[Service]
User=root
Restart=on-failure
RestartSec=5
WorkingDirectory=/data/crud-sample/
ExecStart=java -cp app/crud-sample.jar:conf/ -Dspring.profiles.active=prod org.springframework.boot.loader.JarLauncher
TimeoutStopSec=10
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
