# src/assembly/bin/bard-crud-sample.service
[Unit]
Description=Bard CRUD Sample
Requires=network.target

[Service]
User=root
Restart=on-failure
RestartSec=5
WorkingDirectory=/opt/bard/
ExecStart=java -cp apps/crud-sample/*:conf/crud-sample -Dspring.profiles.active=prod org.springframework.boot.loader.JarLauncher
TimeoutStopSec=10
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
